# Front Row Lion/Data/CoreServices/Front Row.app/Contents/Resources/com.apple.RemoteUI.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>com.apple.RemoteUI</string>
        <key>ProgramArguments</key>
        <array>
                <string>/System/Library/CoreServices/Front Row.app/Contents/MacOS/Front Row</string>
        </array>
        <key>MachServices</key>
        <dict>
                <key>com.apple.RemoteUI</key>
                <true/>
        </dict>
</dict>
</plist>
